# system/airdancer.target
[Unit]
Description=Airdancer

[Install]
WantedBy=default.target
